# net/mqttbroker/files/wsmqttintegrator.init
#!/bin/sh /etc/rc.common

START=50
STOP=90

USE_PROCD=1

start_service() {
        procd_open_instance
        procd_set_param command /usr/bin/wsmqttintegrator
        procd_append_param command legacy socket --retry --retry-base 1 --reconnect
        procd_set_param respawn
        procd_close_instance
}
